# salt/djangowebserver/gunicorn_supervisor.conf
[program:gunicorn]
command = {{ virtualenv }}/bin/gunicorn -c {{ project }}/gunicorn.conf.py
user = {{ user }}
autostart = True
autorestart = True
redirect_stderr = True
