# salt/djangowebserver/gunicorn_supervisor.conf
[program:{{ program_name }}]
command = {{ virtualenv }}/bin/gunicorn -c {{ project }}/gunicorn.conf.py
user = {{ gunicorn_user }}
autostart = True
autorestart = True
redirect_stderr = True
